# vps-backup/koreader-sync.service
[Unit]
Description=KOReader Notes Auto-Sync Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash /root/scripts/sync_notes.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
